# templates/database-backup.timer.j2
[Unit]
Description=Run database backup daily

[Timer]
OnCalendar={{ database_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target